# cgar_collect.service
[Unit]
Description=Collect cgroup data (cgar)

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/cgar_collect