# cgar_collect.timer
[Unit]
Description=Collect cgroup data (cgar)

[Timer]
OnCalendar=*:0/1
Unit=cgar_collect.service

[Install]
WantedBy=multi-user.target